# cookbooks/installers/prometheus/node-exporter/files/service/node-exporter/run
#!/bin/bash

set -e

exec node_exporter --path.rootfs=/host