# cookbooks/installers/prometheus/prometheus/files/service/prometheus/run
#!/bin/bash

set -e

cd /etc/prometheus/

exec prometheus --config.file=/etc/prometheus/prometheus.yml --web.external-url="https://prometheus.local" --web.route-prefix="/"