# cookbook/install/vault/files/service/vault/run
#!/usr/bin/env bash

set -e

exec vault server -config=/etc/vault/config/config.hcl