# cookbooks/installers/prometheus/cadvisor/files/service/cadvisor/run
#!/bin/bash

set -e

exec cadvisor -logtostderr