# cookbooks/installers/prometheus/files/service/prometheus/run
#!/bin/bash

set -e

cd /etc/prometheus/

exec prometheus --config.file=/etc/prometheus/prometheus.yml