# cookbooks/installers/blackbox/files/service/blackbox/run
#!/bin/bash

set -e

exec blackbox_exporter --config.file=/etc/blackbox/blackbox.yml