# install/nginx/files/service/nginx/run
#!/usr/bin/env bash

set -e

pkill nginx || true

mkdir -p /var/log/nginx
chown -R www-data:www-data /var/log/nginx

exec nginx -g "daemon off; master_process on;"