# cookbooks/installers/prometheus/mysql-exporter/service/mysql-exporter/run
#!/bin/bash

set -e

exec mysqld_exporter --config.my-cnf=/etc/mysql-exporter/mysql-exporter.cnf