# cookbooks/installers/redis/files/service/redis/run
#!/bin/bash

set -e

exec redis-server