# cookbooks/installers/mysql/files/service/mysql/run
#!/bin/bash

set -e

exec mysqld_safe --disable-log-bin --sql-mode="NO_ENGINE_SUBSTITUTION"