# cookbook/install/supervisor/files/service/supervisord/run
#!/usr/bin/env bash

set -e

exec supervisord --configuration /etc/supervisor/supervisord.conf --nodaemon