# cookbooks/installers/traefik/files/service/traefik/run
#!/usr/bin/env bash

set -e

exec traefik --configfile="/etc/traefik/<TRAEFIK_CONFIG>.yml"