# cookbooks/installers/alert-manager/files/service/alert-manager/run
#!/bin/bash

set -e

exec alertmanager --config.file=/etc/alert-manager/alertmanager.yml