# cookbooks/installers/memcache/files/service/memcached-sessions/run
#!/usr/bin/env bash

set -e

exec s6-setuidgid memcache /usr/bin/memcached -m 64 -p 11210 -u memcache