# cookbooks/installers/memcache/files/service/memcached/run
#!/usr/bin/env bash

set -e

exec s6-setuidgid memcache /usr/bin/memcached -m 256 -p 11211 -u memcache